# files/ferengi@.timer
[Unit]
Description=Daily weather data update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
