# files/ferengi@.service
[Unit]
Description=Update FERENGI data
After=network.target

[Service]
User=ferengi
Group=ferengi
Type=oneshot
ExecStart=/usr/local/bin/ferengi %i
