# files/ferengi@.service
[Unit]
Description=Update weather data
After=network.target

[Service]
User=ferengi
Group=ferengi
Type=oneshot
ExecStart=/usr/bin/ferengi weather
